# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
autostart=true
autorestart=true
command=/etc/init.d/elasticsearch start
